# packaging/systemd/sc-metrics-agent-updater.service
[Unit]
Description=SC Metrics Agent Updater
Documentation=https://github.com/strettch/sc-metrics-agent
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/sc-metrics-agent-updater.sh
User=root
Group=root
StandardOutput=journal
StandardError=journal
RemainAfterExit=no

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lock /var/lib/apt /var/cache/apt /etc/apt

# Timeout settings
TimeoutStartSec=300
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target